# Pinger/Resources/renderphinepinger.service
[Unit]
Description=Tracks and restarts renderphine if frozen

[Service]
Type=notify
ExecStart=/usr/sbin/pinger

[Install]
WantedBy=multi-user.target